# init.d/wifimanager
#!/bin/sh /etc/rc.common

START=79

BIN=/root/test_lim/wifi_manager/test_wifimanager.sh

boot() {
    echo "wifi_manager boot"
    $BIN -boot
}
