# init.d/wifimanager
#!/bin/sh /etc/rc.common

START=79

BIN=/usr/bin/wifimanager

boot() {
    echo "wifimanager boot"
    $BIN -boot
}
